# backend/deploy/eznihongo-api.service
[Unit]
Description=EzNihongo API (Node.js)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/eznihongo/backend
EnvironmentFile=/var/www/eznihongo/backend/.env
ExecStart=/usr/bin/node /var/www/eznihongo/backend/src/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eznihongo-api

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/eznihongo/backend /var/www/eznihongo/uploads
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
